# files/openhab.conf
[program:openhab]
command=/opt/openhab/start.sh
directory=/opt/openhab
stopwaitsecs=30
