# files/openhab_debug.conf
[program:openhab_debug]
command=/opt/openhab/start_debug.sh
directory=/opt/openhab
stopwaitsecs=30
autostart=false
